# hosts/cerebellum/systemd/psyched-gps.service
[Unit]
Description=Psyched gps Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pete/psyched
Environment=HOST=cerebellum
ExecStart=/home/pete/psyched/tools/systemd_entrypoint.sh /home/pete/psyched/hosts/cerebellum/modules/gps/launch.sh
ExecStop=/home/pete/psyched/tools/systemd_entrypoint.sh /home/pete/psyched/hosts/cerebellum/modules/gps/shutdown.sh
Restart=on-failure
RestartSec=2
User=root
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
