# hosts/cerebellum/systemd/psyched-ear.service
[Unit]
Description=Psyched ear Service
# Ensure networking is fully configured before starting ROS2 nodes that may
# rely on network interfaces for DDS discovery.
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pete/psyched
Environment=HOST=cerebellum
# ExecStart: run the launch script via the entrypoint which will source
# the ROS environment and then exec the launch script. Keep it simple.
ExecStart=/home/pete/psyched/tools/systemd_entrypoint.sh /home/pete/psyched/hosts/cerebellum/modules/ear/launch.sh
# ExecStop: some shutdown scripts may not be marked executable; run the
# entrypoint which will source env and then invoke the shutdown script via
# an explicit shell command to avoid permission-denied exit code 126.
# We use bash -c to ensure the shutdown script is executed by the shell.
ExecStop=/home/pete/psyched/tools/systemd_entrypoint.sh bash -c "\"/home/pete/psyched/hosts/cerebellum/modules/ear/shutdown.sh\""
Restart=on-failure
RestartSec=2
User=pete
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
