# hosts/cerebellum/systemd/psyched-pilot.service
# Generated by psh systemd
[Unit]
Description=Psyched Pilot UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pete/psyched
Environment=HOST=cerebellum
Environment=PSH_MODULE_NAME=pilot
Environment=PSH_MODULE_CONFIG=/home/pete/psyched/hosts/cerebellum/config/pilot.yaml
ExecStart=/home/pete/psyched/tools/systemd_entrypoint.sh /home/pete/psyched/modules/pilot/pilot.launch.sh
ExecStop=/home/pete/psyched/tools/systemd_entrypoint.sh /home/pete/psyched/modules/pilot/pilot.shutdown.sh
Restart=on-failure
RestartSec=4
User=pete
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
