# hosts/cerebellum/systemd/psyched-chat.service
# Generated by psh systemd
[Unit]
Description=Psyched Chat Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pete/psyched
Environment=HOST=cerebellum
Environment=PSH_MODULE_NAME=chat
Environment=PSH_MODULE_CONFIG=/home/pete/psyched/hosts/cerebellum/config/chat.yaml
ExecStart=/home/pete/psyched/tools/systemd_entrypoint.sh bash -lc "set -euo pipefail\nSYSTEM_PROMPT_VAL=\"${CHAT_SYSTEM_PROMPT:-${SYSTEM_PROMPT:-You_are_a_helpful_assistant._Reply_in_one_concise_sentence.}}\"\nCONV_TOPIC_VAL=\"${CHAT_CONVERSATION_TOPIC:-${CONVERSATION_TOPIC:-/conversation}}\"\nVOICE_TOPIC_VAL=\"${CHAT_VOICE_TOPIC:-${VOICE_TOPIC:-/voice}}\"\nMODEL_VAL=\"${CHAT_MODEL:-${OLLAMA_MODEL:-tinyllama}}\"\nOLLAMA_HOST_VAL=\"${OLLAMA_HOST:-http://localhost:11434}\"\nMAX_HISTORY_VAL=\"${CHAT_MAX_HISTORY:-${MAX_HISTORY:-20}}\"\nexec ros2 launch chat chat.launch.py system_prompt:=\"${SYSTEM_PROMPT_VAL}\" conversation_topic:=\"${CONV_TOPIC_VAL}\" voice_topic:=\"${VOICE_TOPIC_VAL}\" model:=\"${MODEL_VAL}\" ollama_host:=\"${OLLAMA_HOST_VAL}\" max_history:=\"${MAX_HISTORY_VAL}\""
ExecStop=/home/pete/psyched/tools/systemd_entrypoint.sh bash -lc "set -euo pipefail\nPATTERN=\"ros2 launch chat chat.launch.py\"\nTIMEOUT=${TIMEOUT:-10}\nif ! pgrep -f \"$PATTERN\" >/dev/null 2>&1; then\n  echo \"[chat/shutdown] No matching processes found for pattern: $PATTERN\"\n  exit 0\nfi\n\necho \"[chat/shutdown] Sending SIGTERM to processes for: $PATTERN\"\npkill -TERM -f \"$PATTERN\" || true\n\nfor ((i=0; i<TIMEOUT; i++)); do\n  sleep 1\n  if ! pgrep -f \"$PATTERN\" >/dev/null 2>&1; then\n    echo \"[chat/shutdown] All processes stopped\"\n    exit 0\n  fi\ndone\n\necho \"[chat/shutdown] Forcing SIGKILL for remaining processes\"\npkill -KILL -f \"$PATTERN\" || true"
Restart=on-failure
RestartSec=4
User=pete
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
